# roles/proxmox_nginx/templates/nginx.service.j2
{{ ansible_managed | comment }}

[Unit]
Requires=pve-cluster.service
After=pve-cluster.service